# backup.service
[Unit]
Description=Daily Backup

[Service]
Type=simple
ExecStart=/opt/backup.sh
