# manual_backup.service
[Unit]
Description=Manual Backup

[Service]
Type=simple
ExecStart=/opt/manual_backup.sh
